# v5-2/web/supervisord.conf
[supervisord]
loglevel = info
logfile=/dev/null
logfile_maxbytes=0
pidfile = /tmp/supervisord.pid
umask = 022
nodaemon=true
user=root

;logfile = /web/logs/supervisord.log
;logfile_maxbytes = 50MB
;logfile_backups = 30
;loglevel = info

[supervisorctl]

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port = *:9001
username=admin
password=d5a122d9

[program:data_science]
numprocs = %(ENV_CORES)s
numprocs_start = 1
process_name = 100%(process_num)s

command=python3 /web/share_global.py --port=100%(process_num)s

stopwaitsecs=5
startretries=5
autostart=true
autorestart=true

stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
